# config/recoverymanager.service.in
###############################################################################
#
#  Coredump manager
#
###############################################################################
[Unit]
Description=Crasmanager service

[Service]
Type=simple
User=root
Group=root
ExecStart=@install_prefix@/bin/crashmanager
RuntimeDirectory="crashmanager"
StateDirectory="crashmanager"
WatchdogSec=30

[Install]
WantedBy=multi-user.target
