# config/recoverymanager.service.in
###############################################################################
#
#  Recovery manager
#
###############################################################################
[Unit]
Description=Recovery Manager service

[Service]
Type=simple
User=root
Group=root
ExecStart=@install_prefix@/bin/recoverymanager
RuntimeDirectory="rmgr"
StateDirectory="rmgr"
WatchdogSec=30

[Install]
WantedBy=multi-user.target
