# config/recoverymanager.service.in
###############################################################################
#
#  RecoveryManager service
#
###############################################################################
[Unit]
Description=Recovery Manager service
Wants=dbus.service
After=dbus.service

[Service]
Type=notify
User=root
Group=root
ExecStart=@install_prefix@/bin/recoverymanager
RuntimeDirectory="recoverymanager"
StateDirectory="recoverymanager"
WatchdogSec=30

[Install]
WantedBy=multi-user.target
