# testing/echoservice/echos4.service
[Unit]
Description=Echo service S4

[Service]
Type=simple
ExecStart=/usr/bin/echoservice S4

[Install]
WantedBy=multi-user.target
