# testing/echoservice/echos5.service
[Unit]
Description=Echo service S5

[Service]
Type=simple
ExecStart=/usr/bin/echoservice S5

[Install]
WantedBy=multi-user.target
