# testing/echoservice/echos9.service
[Unit]
Description=Echo service S9

[Service]
Type=simple
ExecStart=/usr/bin/echoservice S9

[Install]
WantedBy=multi-user.target
