# testing/echoservice/echos6.service
[Unit]
Description=Echo service S6

[Service]
Type=simple
ExecStart=/usr/bin/echoservice S6

[Install]
WantedBy=multi-user.target
